# resources/postgres.cron
* * * * * /usr/bin/echo "insert into monitoring (value) select floor(random() * 100)::int;" | /usr/bin/psql monitoring
